# server/server_configs/on_boot.service
[Service]
ExecStart=/root/boot.sh
[Install]
WantedBy=default.target
